# webp2jpg.service
[Unit]
Description=WebP to JPG Converter Service
After=network.target

[Service]
Type=simple
User=ovc
WorkingDirectory=/home/ovc/RustroverProjects/webp2jpg
Environment="WEBP2JPG_WATCH_DIR=/home/ovc/Bilder/_Speakers"
Environment="WEBP2JPG_CROP=true"
ExecStart=/home/ovc/RustroverProjects/webp2jpg/target/release/webp2jpg
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target